# rootfs/etc/services.d/shelly-scanner/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Take down the S6 supervision tree when the service fails
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

declare exit_code
readonly exit_code_container=$(</run/s6-linux-init-container-results/exitcode)
readonly exit_code_service="${1}"
readonly exit_code_signal="${2}"
readonly service="Shelly Scanner"

bashio::log.info "Service ${service} exited with code ${exit_code_service} (by signal ${exit_code_signal})"

if [[ "${exit_code_service}" -eq 256 ]]; then
  if [[ "${exit_code_container}" -eq 0 ]]; then
    bashio::log.info "Requesting Home Assistant restart..."
    bashio::homeassistant.api.request "core/restart"
  else
    bashio::log.warning "Service restart requested, but Home Assistant is already stopping."
  fi
else
  bashio::log.info "Halting Home Assistant add-on..."
fi
